# ocr-ci2/deploy/prod/native/ocr-gateway.service.example
# systemd unit for production native OCR Gateway (Linux)
#
# Install:
#   sudo cp deploy/prod/native/ocr-gateway.service.example /etc/systemd/system/ocr-gateway.service
#   sudo sed -i "s|/opt/ocr-ci2|$(pwd)|g" /etc/systemd/system/ocr-gateway.service
#   sudo sed -i "s|YOUR_USER|$USER|g" /etc/systemd/system/ocr-gateway.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now ocr-gateway
#
# Requires: deploy/prod/native/gateway.env

[Unit]
Description=OCR Gateway (OpenCodeReview MR reviews)
After=network.target

[Service]
Type=simple
User=YOUR_USER
WorkingDirectory=/opt/ocr-ci2
Environment=PYTHONPATH=/opt/ocr-ci2
EnvironmentFile=/opt/ocr-ci2/deploy/prod/native/gateway.env
ExecStart=/usr/bin/python3 -m uvicorn gateway.main:app --host 0.0.0.0 --port 8010
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
